# etc/systemd/system/systembackup.service
[Unit]
Description=create a system backup, tar, and copy to HDD

[Service]
ExecStart=/usr/local/sbin/systembackup.sh
